# scripts/services/clock_sync/run
#!/bin/sh
#
# clock_sync - set the system clock accurately from an HTTP Date header as soon as any network path comes up

exec 2>&1

SAVED=/var/local/lib/clock_saved
SYNCED_FLAG=/run/clock_synced

# Captive-portal style endpoints first, then plain hosts whose headers still carry Date. 
# tried in order until one answers
URLS="http://www.gstatic.com/generate_204 http://connectivitycheck.gstatic.com/generate_204 http://detectportal.firefox.com/success.txt http://1.1.1.1"

# Reformat an RFC 1123 HTTP Date into the "YYYY-MM-DD HH:MM:SS" form that busybox `date -s` reliably accepts
http_date_to_settable() {
  echo "$1" | awk '
    { split("Jan Feb Mar Apr May Jun Jul Aug Sep Oct Nov Dec", names, " ")
      for (i = 1; i <= 12; i++) mon[names[i]] = sprintf("%02d", i)
      # Fields: 1=Day,  2=DD  3=Mon  4=YYYY  5=HH:MM:SS  6=GMT
      if (NF >= 5 && mon[$3] != "") printf "%s-%s-%s %s\n", $4, mon[$3], $2, $5 }'
}

set_from_http() {
  for url in $URLS; do
    # -sI: silent head request. -m 5: 5s timeout. Strip CRs, pull the Date header value (case-insensitive)
    hdr=$(curl -sI -m 5 "$url" 2>/dev/null | tr -d '\r' \
      | awk -F': ' 'tolower($1)=="date"{ sub(/^[^:]*: /,""); print; exit }')
    [ -z "$hdr" ] && continue

    settable=$(http_date_to_settable "$hdr")
    [ -z "$settable" ] && { echo "clock_sync: unparseable Date '${hdr}' from ${url}"; continue; }

    # -u: the HTTP Date is UTC, so set the clock in UTC.
    if date -u -s "$settable" > /dev/null 2>&1; then
      echo "clock_sync: clock set from ${url} -> $(date -u '+%Y-%m-%d %H:%M:%S UTC')"
      date +%s > "$SAVED" 2>/dev/null
      return 0
    fi
    echo "clock_sync: 'date -s ${settable}' failed (from ${url})"
  done
  return 1
}

# Initial fast sync: retry every 3s until a network path is up
until set_from_http; do
  sleep 3
done

# Signal that wall time is now trustworthy
# here if we later want the daemon to hold playback position until clock is accurate 
touch "$SYNCED_FLAG" 2>/dev/null

# re-sync hourly to correct drift
while true; do
  sleep 3600
  set_from_http > /dev/null 2>&1 || echo "clock_sync: periodic re-sync failed (network down?), will retry next hour"
done
